# biomoqa-rag-v3-fast.service
[Unit]
Description=BioMoQA RAG V3.1 Ultra-Fast API Server
After=network.target

[Service]
Type=simple
User=egaillac
WorkingDirectory=/home/egaillac/BioMoQA-RAG
Environment="PATH=/home/egaillac/BioMoQA-RAG/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/egaillac/BioMoQA-RAG/venv/bin/python3 -m uvicorn api_server_v3_fast:app --host 0.0.0.0 --port 9000
Restart=always
RestartSec=10
StandardOutput=append:/home/egaillac/BioMoQA-RAG/logs/api_v3_fast.log
StandardError=append:/home/egaillac/BioMoQA-RAG/logs/api_v3_fast_error.log

[Install]
WantedBy=multi-user.target
